# system/init.d/greetd
#!/sbin/openrc-run
# OpenRC service for greetd.
#
# The Gentoo gui-libs/greetd package ships ONLY /usr/lib/systemd/system/greetd.service
# and no init script, so on an OpenRC system `rc-update add greetd default` fails
# with nothing to add. This file fills that gap.
description="greetd login manager"
supervisor=supervise-daemon
command="/usr/bin/greetd"
command_args="--config /etc/greetd/config.toml"
respawn_delay=5
respawn_max=0

depend() {
    need dbus
    after elogind bootmisc
    provide display-manager
}
